# ops/launchd/io.github.deepseek-harness-corpus-supervisor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.deepseek-harness.corpus-supervisor</string>
  <key>ProgramArguments</key>
  <array>
    <string>__DEEPSEEK_HARNESS_BIN__</string>
    <string>corpus</string>
    <string>supervise</string>
    <string>--max-cycles</string>
    <string>10000</string>
    <string>--interval-ms</string>
    <string>5000</string>
    <string>--max-jobs-per-cycle</string>
    <string>4</string>
    <string>--max-iterations-per-job</string>
    <string>10</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DEEPSEEK_HARNESS_ARTIFACT_DIR</key>
    <string>__DEEPSEEK_HARNESS_ARTIFACT_DIR__</string>
    <key>DEEPSEEK_HARNESS_STATE_DIR</key>
    <string>__DEEPSEEK_HARNESS_STATE_DIR__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>__DEEPSEEK_HARNESS_LOG_DIR__/deepseek-harness-corpus-supervisor.log</string>
  <key>StandardErrorPath</key>
  <string>__DEEPSEEK_HARNESS_LOG_DIR__/deepseek-harness-corpus-supervisor.error.log</string>
</dict>
</plist>
